# infrastructure/systemd/lunchlineup-retention-purge.timer
[Unit]
Description=Daily LunchLineup 30-day application-data purge

[Timer]
OnCalendar=*-*-* 03:17:00
RandomizedDelaySec=30m
Persistent=true
Unit=lunchlineup-retention-purge.service

[Install]
WantedBy=timers.target
