# infrastructure/systemd/lunchlineup-backup.timer
[Unit]
Description=Daily LunchLineup encrypted offsite database backup

[Timer]
OnCalendar=*-*-* 02:17:00
RandomizedDelaySec=30m
Persistent=true
Unit=lunchlineup-backup.service

[Install]
WantedBy=timers.target
