# infrastructure/systemd/lunchlineup-public-web-probe.service
[Unit]
Description=LunchLineup bounded public web availability probe
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=lunchlineup
Group=lunchlineup
EnvironmentFile=/etc/lunchlineup/public-web-probe.env
ExecStart=/usr/bin/bash /opt/lunchlineup/infrastructure/control/public-web-probe.sh
TimeoutStartSec=25s
UMask=0027
NoNewPrivileges=true
PrivateDevices=true
PrivateTmp=true
ProtectControlGroups=true
ProtectHome=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectSystem=strict
ReadOnlyPaths=/opt/lunchlineup /etc/lunchlineup
ReadWritePaths=/var/lib/node_exporter/textfile_collector
RestrictAddressFamilies=AF_INET AF_INET6
RestrictRealtime=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native
SystemCallFilter=@system-service
